# nomad/templates/nomad.service.j2
[Unit]
Description=Nomad Agent

[Service]
ExecStart=/usr/sbin/nomad agent -config /etc/nomad.conf

[Install]
WantedBy=multi-user.target
